# conf.d/fcserver-partyhat.conf
[program:fcserver-partyhat]
command=/opt/fadecandy/bin/fcserver-rpi /usr/src/app/fcserver-partyhat.json
autostart=false
