# conf.d/fcserver-partyhat.conf
[program:fcserver]
command=/opt/fadecandy/bin/fcserver-rpi /app/fcserver-partyhat.json
