# conf.d/fcserver.conf
[program:fcserver]
command=/opt/fadecandy/bin/fcserver-rpi /app/fcserver.json
;command=/opt/fadecandy/bin/fcserver-rpi %(ENV_FADECANDY_CONFIG)s
environment=FADECANDY_CONFIG="/app/fcserver.json"
