# conf.d/fcserver.conf
[program:fcserver]
;command=/opt/fadecandy/bin/fcserver-rpi /usr/src/app/fcserver.json
command=/opt/fadecandy/bin/fcserver-rpi %(ENV_FADECANDY_CONFIG)s
