# conf.d/partyhat-rainbow.conf
[program:partyhat-rainbow]
command=python %(ENV_PARTYHAT_SCRIPT)s
environment=PARTYHAT_SCRIPT="/opt/fadecandy/examples/python/chase.py",FADECANDY_URL="localhost:7891"
