# conf.d/partyhat-rainbow.conf
[program:partyhat-rainbow]
command=python %(ENV_PARTYHAT_SCRIPT)s
environment=PARTYHAT_SCRIPT="/usr/src/app/partyhat/rainbow.py",
